# debian/prerm
#!/bin/sh
systemctl disable wb-knxd-config.service
